# systemd/symphonika.service
[Unit]
Description=Symphonika orchestrator daemon
After=graphical-session.target network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
Restart=on-failure
RestartSec=5s

# Inherit a sane PATH so spawned providers (claude, codex) and tools
# (gh, esbmc, cargo, ...) resolve under user-systemd.
Environment=PATH=%h/.npm-global/bin:%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin

# Resolve GITHUB_TOKEN from `gh auth token` at each (re)start so this
# survives token rotation. Fails closed if gh returns empty so the
# daemon never starts without a token. `exec` replaces the shell so
# the node process becomes the service's Main PID.
#
# Customise the ExecStart path to match your install:
#   * npm global       -> %h/.npm-global/bin/symphonika daemon
#   * source checkout  -> /usr/bin/node /path/to/symphonika/dist/cli.js daemon
ExecStart=/bin/sh -c 't=$(gh auth token); [ -n "$t" ] || { echo "ERROR: gh auth token returned empty"; exit 1; }; export GITHUB_TOKEN="$t"; exec %h/.npm-global/bin/symphonika daemon'

# Keep the daemon and everything it spawns in its own cgroup slice.
# A scope-wide OOM in a spawned tool (compiler, verifier, ...) no
# longer tears down whichever terminal scope you happened to launch
# the daemon from.
Slice=symphonika.slice

# Journald sees stdout/stderr. View with:
#   journalctl --user -u symphonika -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
